# Library/org.zfs.snapshot.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>org.zfs.snapshot.daily</string>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/sbin/zfs-auto-snapshot.sh</string>
                <string>--syslog</string>
                <string>--quiet</string>
                <string>--keep=31</string>
                <string>--label=daily</string>
                <string>//</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>StartInterval</key>
        <integer>86400</integer>
</dict>
</plist>
